# wol-off.service
[Unit]
Description=Disable Wake On Lan
 
[Service]
Type=oneshot
ExecStart = /usr/sbin/ethtool -s eno1 wol d
 
[Install]
WantedBy=basic.target
